# tests/03-systemd-inside/fs/etc/systemd/system/sleep.service
[Unit]
Description=sleep forever

[Service]
Type=exec
ExecStart=/bin/bash -c "sleep infinity"

[Install]
WantedBy=default.target
