# staff/services/containers-ensure-health.timer
[Unit]
Description=schedule ensure all containers status is healthy every 5min

[Timer]
OnBootSec=15min
OnUnitInactiveSec=5min
