# staff/services/services@root.target
[Unit]
Description=System Worker Services
Requires=services-pre.target
After=services-pre.target
AllowIsolate=no
