# staff/services/services@user.target
[Unit]
Description=System Worker Services
Requires=default.target services-pull-all-images.service
After=default.target services-pull-all-images.service
AllowIsolate=no
