# staff/systemd-filesystem/basic/fs/usr/local/lib/systemd/system/container-getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noclear --login-program /entrypoint/login.sh --autologin root --keep-baud pts/%I 115200,38400,9600 $TERM
Restart=always
