# staff/systemd-filesystem/basic/fs/usr/local/lib/systemd/system/console-getty.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noclear --login-program /entrypoint/login.sh --autologin root --keep-baud console 115200,38400,9600 $TERM
ExecStopPost=/usr/bin/systemctl poweroff
StandardInput=tty
StandardOutput=tty

[Install]
WantedBy=multi-user.target
