# services/create-dnsmasq-config.service
[Unit]
After=systemd-resolved.service
[Service]
ExecStartPre=!/usr/bin/env bash "__SCRIPT__"
