# services/create-dnsmasq-config.service
[Unit]
After=systemd-resolved.service libvirt-guests.service libvirtd.service libvirtd.socket
[Service]
ExecStartPre=!/usr/bin/env bash "__SCRIPT__"
