# services/services-spin-up.service
[Unit]
Description=Delay 30s on boot

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash "__SCRIPT__"
Restart=no
RemainAfterExit=yes
