# staff/systemd-filesystem/basic/fs/etc/profile.d/80-env.sh
set -a
source /etc/environment
source /run/.userenvironments
set +a
