# staff/systemd-filesystem/basic/fs/etc/profile.d/80-env.sh
set -a
source /run/.userenvironments
source /etc/environment
set +a
