# staff/systemd-filesystem/nginx_attach/fs/usr/lib/systemd/system/nginx-attach-retry.timer
[Unit]
Description=retry nginx-attach.service when it fails after 5min
RefuseManualStart=yes

[Timer]
OnActiveSec=5min
RemainAfterElapse=no
